# utils/skel/.config/autostart/firstlogin.desktop
[Desktop Entry]
Type=Application
Name=First Login Welcome Program
Exec=/usr/bin/firstlogin
