# utils/etc/profile.d/extrapaths.sh
pathprepend /usr/local/lib/pkgconfig PKG_CONFIG_PATH
pathprepend /usr/local/share/pkgconfig PKG_CONFIG_PATH
pathprepend /usr/local/sbin
pathprepend /usr/local/bin
pathappend /usr/share/man  MANPATH
pathappend /usr/share/info INFOPATH
